# docker/gunicorn.conf

[program:gunicorn]
environment=FLASK_CONFIG=development
command=/usr/bin/gunicorn manage:app -b localhost:5000
directory=/deploy
syslog_addr=unix:///dev/log#dgram
syslog